# configs/systemd/incognito-firewall.service
[Unit]
Description=Incognito OS baseline firewall
Before=network-pre.target
Wants=network-pre.target

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/incognito-firewall-base.sh
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
